# systemd/expire.timer
[Timer]
OnCalendar=22:45

[Install]
WantedBy=timers.target
